# flask_dokerfile/supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /root/uwsgi.ini